# setup/systemd/viam-server.service
[Unit]
Description=Viam Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/viam-server --config /etc/viam.json
Restart=on-failure
RestartSec=5
WatchdogSec=30

[Install]
WantedBy=multi-user.target
