# scripts/ironsight-touch.service
[Unit]
Description=IronSight Interactive Touch Display
After=viam-server.service
Wants=viam-server.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/andrew/Viam-Staubli-Apera-PLC-Mobile-POC/scripts/ironsight-touch.py
Restart=on-failure
RestartSec=5
User=root
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/ironsight/env
# Needs root for /dev/fb0, /dev/input/event*, and systemctl commands
# ANTHROPIC_API_KEY loaded from /etc/ironsight/env

[Install]
WantedBy=multi-user.target
